# microKSBScanner/build/microKSBScanner.service
[Unit]
Description=microKSBService

[Service]
WorkingDirectory=/opt/microKSBService
ExecStart=/opt/microKSBService/microKSBService
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=microKSBScanner

[Install]
WantedBy=multi-user.target